# files/default/runit/agent/run
#!/bin/bash

set -e

export PATH=/var/vcap/bosh/bin:$PATH
exec 2>&1
exec /var/vcap/bosh/bin/bosh_agent -c -I $(cat /etc/infrastructure)
