# files/default/runit/agent/run
#!/bin/bash

set -e

export PATH=/var/vcap/bosh/bin:$PATH
exec 2>&1
exec /var/vcap/bosh/bin/bosh_agent -c -P microcloud -I $(cat /etc/infrastructure) -n https://vcap:vcap@0.0.0.0:6969
